# deploy/systemd/boltz-worker-cpu.service
[Unit]
Description=V-Bio CPU Worker (docker compose)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/data/V-Bio/deploy/docker
ExecStart=/usr/bin/docker compose -f DOCKER_STACK_WORKER_CPU.compose.yml --env-file DOCKER_STACK_WORKER_CPU.env up -d --build
ExecStop=/usr/bin/docker compose -f DOCKER_STACK_WORKER_CPU.compose.yml --env-file DOCKER_STACK_WORKER_CPU.env down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
