# deploy/systemd/boltz-central.service
[Unit]
Description=V-Bio Central Stack (docker compose)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/data/V-Bio/deploy/docker
ExecStart=/usr/bin/docker compose -f DOCKER_STACK_CENTRAL.compose.yml --env-file DOCKER_STACK_CENTRAL.env up -d --build
ExecStop=/usr/bin/docker compose -f DOCKER_STACK_CENTRAL.compose.yml --env-file DOCKER_STACK_CENTRAL.env down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
